# src/fish/weatherService/weather-fetch.service
[Unit]
Description=Fetch weather data
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/youruser/scripts/fetch_weather.sh
