# src/fish/weatherService/weather-fetch.service
[Unit]
Description=Fetch weather data
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/fetch_weather.sh

[Install]
WantedBy=multi-user.target
